# debian/postinst
#!/bin/sh
# postinst script for librusec
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

case "$1" in
    configure)
    	for i in 0 1 2 3 4 5 6 7 8 9 a b c d e f; do
		for j in 0 1 2 3 4 5 6 7 8 9 a b c d e f; do
			mkdir -p /usr/share/drupal6/b/$i$j;
		done
	done
	chown www-data:www-data /usr/share/drupal6/b/??
	db_get librusec/mysql_root_password || true
	pwd=$RET
	if ! echo SHOW TABLES | mysql -uroot -p"$pwd" drupal6 | grep -qw librusec; then
		mysql -uroot -p"$pwd" drupal6 < /usr/share/librusec/schema.sql
	fi
	db_set librusec/mysql_root_password ""
	db_fset librusec/mysql_root_password seen false
	db_go
	db_stop
	ln -sf /usr/share/librusec/apache.conf /etc/apache2/conf.d/librusec.conf
	rm -f /usr/share/drupal6/.htaccess
	ln -s /usr/share/librusec/htaccess /usr/share/drupal6/.htaccess
	ln -sf modules/librusec/AJAX.php /usr/share/drupal6
    	invoke-rc.d apache2 reload
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

